# monitor/deploy/systemd/jetstream-monitor.service
[Unit]
Description=Jetstream Monitor - Compare turbo stream message rates
After=network.target
Documentation=https://github.com/messijo/jetstream-turbo-rs

[Service]
Type=simple
User=jetstream-monitor
Group=jetstream-monitor
WorkingDirectory=/opt/jetstream-monitor
ExecStart=/opt/jetstream-monitor/jetstream-monitor
Restart=always
RestartSec=10

EnvironmentFile=/opt/jetstream-monitor/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=jetstream-monitor

LimitNOFILE=65536
LimitNPROC=4096

PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/opt/jetstream-monitor

[Install]
WantedBy=multi-user.target
