# rust/deploy/systemd/jetstream-turbo.service
[Unit]
Description=Jetstream Turbo - Bluesky firehose processor
After=network.target
Documentation=https://github.com/messijo/jetstream-turbo-rs

[Service]
Type=simple
User=jetstream-turbo
Group=jetstream-turbo
WorkingDirectory=/opt/jetstream-turbo
ExecStart=/opt/jetstream-turbo/jetstream-turbo
Restart=always
RestartSec=10

EnvironmentFile=/opt/jetstream-turbo/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=jetstream-turbo

LimitNOFILE=65536
LimitNPROC=4096

PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/opt/jetstream-turbo

[Install]
WantedBy=multi-user.target
